# etc/systemd/system/ArchiBoot.service
[Unit]
Description=ArchiBoot service

[Service]
ExecStart=/root/autostart.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
